# Group/Windows-Server-2025/ossec.conf
 <agent_config>
    <!-- Shared agent configuration here -->
    <!-- ══════════════════════════════════════════════════════
       TIER 1 — CRITICAL
       ══════════════════════════════════════════════════════ -->

  <!-- Core auth, logon, account management -->
  <localfile>
    <location>Security</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Service installs, crashes, driver loads -->
  <localfile>
    <location>System</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Application errors and warnings -->
  <localfile>
    <location>Application</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- AD object changes, LDAP, replication errors -->
  <localfile>
    <location>Directory Service</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- DNS queries, zone changes, errors -->
  <localfile>
    <location>DNS Server</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- AD replication failures -->
  <localfile>
    <location>DFS Replication</location>
    <log_format>eventchannel</log_format>
  </localfile>


  <!-- ══════════════════════════════════════════════════════
       TIER 2 — HIGH (Attacker TTPs)
       ══════════════════════════════════════════════════════ -->

  <!-- 4103/4104 PowerShell module & script block logging -->
  <localfile>
    <location>Microsoft-Windows-PowerShell/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- PowerShell admin events -->
  <localfile>
    <location>Microsoft-Windows-PowerShell/Admin</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Legacy PS engine events 400/403/600 -->
  <localfile>
    <location>Windows PowerShell</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- GPO apply failures, changes — T1484 -->
  <localfile>
    <location>Microsoft-Windows-GroupPolicy/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Kerberos ticket events AS/TGS — T1558 Kerberoasting -->
  <localfile>
    <location>Microsoft-Windows-Kerberos/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- NTLM auth attempts — T1550.002 Pass-the-Hash -->
  <localfile>
    <location>Microsoft-Windows-NTLM/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Netlogon failures, secure channel — T1078 -->
  <localfile>
    <location>Microsoft-Windows-Security-Netlogon/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Scheduled task creation/run — T1053 Persistence -->
  <localfile>
    <location>Microsoft-Windows-TaskScheduler/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- WMI queries — T1047 Lateral movement -->
  <localfile>
    <location>Microsoft-Windows-WMI-Activity/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Blocked EXE/DLL execution — T1562 -->
  <localfile>
    <location>Microsoft-Windows-AppLocker/EXE and DLL</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Blocked script execution -->
  <localfile>
    <location>Microsoft-Windows-AppLocker/MSI and Script</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Defender detections, quarantine — T1562.001 -->
  <localfile>
    <location>Microsoft-Windows-Windows Defender/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Firewall rule changes — T1562.004 -->
  <localfile>
    <location>Microsoft-Windows-Windows Firewall With Advanced Security/Firewall</location>
    <log_format>eventchannel</log_format>
  </localfile>


  <!-- ══════════════════════════════════════════════════════
       TIER 3 — MEDIUM (Lateral Movement & Remote Access)
       ══════════════════════════════════════════════════════ -->

  <!-- RDP local session tracking — T1021.001 -->
  <localfile>
    <location>Microsoft-Windows-TerminalServices-LocalSessionManager/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Incoming RDP connections -->
  <localfile>
    <location>Microsoft-Windows-TerminalServices-RemoteConnectionManager/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- SMB share access auditing — T1021.002 -->
  <localfile>
    <location>Microsoft-Windows-SMBServer/Audit</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- SMB security events -->
  <localfile>
    <location>Microsoft-Windows-SMBServer/Security</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Client-side SMB auth failures -->
  <localfile>
    <location>Microsoft-Windows-SmbClient/Security</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- UAC elevation events — T1548.002 -->
  <localfile>
    <location>Microsoft-Windows-UAC/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Unsigned driver/binary loads -->
  <localfile>
    <location>Microsoft-Windows-CodeIntegrity/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- LAPS local admin password rotation -->
  <localfile>
    <location>Microsoft-Windows-LAPS/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- DNS record modifications -->
  <localfile>
    <location>Microsoft-Windows-DNSServer/Audit</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- BITS transfers — data exfiltration detection -->
  <localfile>
    <location>Microsoft-Windows-Bits-Client/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- WinRM remote command execution -->
  <localfile>
    <location>Microsoft-Windows-WinRM/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- SSH access (OpenSSH installed on your DC) -->
  <localfile>
    <location>OpenSSH/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>


  <!-- ══════════════════════════════════════════════════════
       TIER 4 — LOW (Compliance / Threat Hunting)
       ══════════════════════════════════════════════════════ -->

  <!-- DPAPI key access — credential theft detection -->
  <localfile>
    <location>Microsoft-Windows-Crypto-DPAPI/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Encrypted file system events -->
  <localfile>
    <location>Microsoft-Windows-EFS/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Profile load anomalies -->
  <localfile>
    <location>Microsoft-Windows-User Profile Service/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Winlogon chain events -->
  <localfile>
    <location>Microsoft-Windows-Winlogon/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- USB / removable device detection -->
  <localfile>
    <location>Microsoft-Windows-Kernel-PnP/Device Management</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Network profile changes -->
  <localfile>
    <location>Microsoft-Windows-NetworkProfile/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Time sync — detect time rollback/tampering attacks -->
  <localfile>
    <location>Microsoft-Windows-Time-Service/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Certificate authority deployment changes -->
  <localfile>
    <location>Microsoft-Windows-CertificateServices-Deployment/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>

  <!-- Certificate lifecycle events -->
  <localfile>
    <location>Microsoft-Windows-CertificateServicesClient-Lifecycle-System/Operational</location>
    <log_format>eventchannel</log_format>
  </localfile>
 </agent_config>
